# raspberry-pi/diy-efis/diy-efis.service
[Unit]
Description=DIY EFIS Display

[Service]
ExecStart=/usr/diy-efis/diy-efis
Restart=on-failure

[Install[
WantedBy=multi-user.target
